# deployment/systemd/tpm-live.service
[Unit]
Description=IrsanAI TPM Live Monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/IrsanAI-TPM
ExecStart=/usr/bin/python3 scripts/tpm_cli.py live --history-csv btc_real_24h.csv --poll-seconds 3600
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
